# remote-backup.timer
[Unit]
Description=Run remote-backup every 24h at 12:00pm

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target